# deploy/linglong-mcp.service
[Unit]
Description=Linglong MCP Server
After=network.target

[Service]
Type=simple
User=linglong
WorkingDirectory=/opt/linglong
ExecStart=/opt/linglong/.venv/bin/python -m linglong.mcp
Restart=always
RestartSec=5

# 环境变量
Environment=LL_MCP_AUTH_TOKEN=your-secret-token
Environment=ZHIPU_API_KEY=your-zhipu-key
Environment=SEARXNG_API_KEY=your-searxng-key
Environment=RSSHUB_ACCESS_KEY=your-rsshub-key

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=linglong-mcp

[Install]
WantedBy=multi-user.target
